# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:php5fpm]
command=php5-fpm --fpm-config /etc/php5/fpm/php-fpm.conf -F
